# etc/rc.d/aprsd
#!/bin/ksh
#
# $OpenBSD$

daemon="/usr/sbin/aprsd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
